# deployment/supervisor.template.conf
[program:SITENAME]
command=SITE_FOLDER/bin/gunicorn_start
user=USER
redirect_stderr=true
stdout_logfile=SITE_FOLDER/tmp/logs/SITENAME.log
enviroment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
